# grafstation/configs/grafana/apply.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.lkj.life.dashboard.grafstation.grafana.apply</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TF_VAR_OPENWEATHER_CITY</key>
        <string><%= ow_city %></string>
        <key>TF_VAR_YOUTUBE_PLAYLIST_ID</key>
        <string><%= yt_playlist %></string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/<%= user %>/ghq/github.com/legnoh/life-dashboard/grafana/apply.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/<%= user %>/ghq/github.com/legnoh/life-dashboard/grafana</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/grafana-apply.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/grafana-apply.err</string>
</dict>
</plist>
